# docker/supervisor/conf.d/uvicorn.ini
[program:uvicorn]
command=poetry run uvicorn program.main:app --proxy-headers --host 0.0.0.0 --port 5000
directory=/app/
autostart=true
autorestart=true
startsecs=30
startretries=3
exitcodes=0,2
stopsignal=QUIT
stopwaitsecs=10
user=root
redirect_stderr=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/uvicorn.stdout.log